# .config/systemd/user/mega-cmd.service
[Unit]
Description=Sync stuff

[Service]
ExecStart=/usr/bin/mega-cmd-server
Restart=always

[Install]
WantedBy=default.target
