# .config/systemd/user/mpris-proxy.service
[Unit]
Description=Make bluetooth headset buttons work

[Service]
ExecStart=/usr/bin/mpris-proxy
Restart=always

[Install]
WantedBy=default.target
