# .config/systemd/user/wsession.target
[Unit]
Description=Wayland session running
BindsTo=graphical-session.target
